# dynamic-scaling/reth-scaling.service
[Unit]
Description=Rust Seed Transaction Runner
After=network.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/scripts/reth-scaling/dynamic-scaling
Environment=DOTENV_PATH=/home/ubuntu/scripts/reth-scaling/dynamic-scaling/.env
Environment=RUST_LOG=info

# Remove the extra space after ExecStart=
ExecStart=/home/ubuntu/scripts/reth-scaling/target/release/seed send-eth-coh-no-receipt --from-node 5 --to-node 4 --num-txs 100 --amount-wei 1 --zero-gas-price --delay-secs 1

# Restart configuration
Restart=on-failure
RestartSec=30
StartLimitBurst=5

# Logging
StandardOutput=append:/var/log/reth-scaling.log
StandardError=append:/var/log/reth-scaling.log

# Resource limits
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target